# config-sync/claude/skills/trading-morning-levels/com.trading.morninglevels.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trading.morninglevels</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/trp/.claude/skills/trading-morning-levels/morning_levels.sh</string>
    </array>

    <!-- Run Mon–Fri at 8:30 AM local time (Weekday: 1=Mon … 5=Fri) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/trp/.claude/skills/trading-morning-levels/morning_levels.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/trp/.claude/skills/trading-morning-levels/morning_levels.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/trp</string>
        <key>DISCORD_LEVELS_WEBHOOK</key>
        <string>https://discord.com/api/webhooks/1503394040664821900/Qe99bEHiV458s8cPMCgSk0ozEG6E5_tSKowusEV8R5Vu4owxmsGz1RGHwzknckSQ_6Xm</string>
    </dict>
</dict>
</plist>
